# gameprj/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
svnserve -d --listen-port 8888 -r /home/isq/svndata
cd /home/isq/gameprj/
nohup ./gserver1.0_watch &
